# templates/servo/sanoid/syncoid.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Syncoid Periodic Timer

[Timer]
OnCalendar=*-*-* 3:37:00
Persistent=true

[Install]
WantedBy=timers.target